# doc-server/debian/DEBIAN/postrm
#!/bin/sh
set -e

/usr/sbin/update-rc.d -f docserver remove

