# doc-server/debian/DEBIAN/postinst
#!/bin/sh
set -e

/usr/sbin/update-rc.d docserver defaults
